# build/docker/supervisord.conf
; supervisord.conf
[unix_http_server]
;file=/tmp/supervisor.sock   ; (the path to the socket file)
;修改为 /var/run 目录，避免被系统删除
file=/app/running/storage/logs/supervisord/supervisor.sock   ; (the path to the socket file)
;chmod=0700                 ; socket file mode (default 0700)
;chown=nobody:nogroup       ; socket file uid:gid owner
;username=user              ; (default is no username (open server))
;password=123               ; (default is no password (open server))

[supervisord]
nodaemon=true
logfile=/app/running/storage/logs/supervisord/supervisord.log
pidfile=/app/running/storage/logs/supervisord/supervisord.pid

[inet_http_server]
port=0.0.0.0:9001
username=admin
password=secretpassword

[program:qcluster]
command=python manage.py qcluster
directory=/app
autostart=true
autorestart=true
stopasgroup=true
stdout_logfile=/app/running/storage/logs/supervisord/qcluster.log
stderr_logfile=/app/running/storage/logs/supervisord/qcluster.err
priority=100

[program:gunicorn]
command=gunicorn -c /app/build/docker/gunicorn.conf.py starwizAi.wsgi:application
directory=/app
autostart=true
autorestart=true
stopasgroup=true
stdout_logfile=/app/running/storage/logs/supervisord/gunicorn.log
stderr_logfile=/app/running/storage/logs/supervisord/gunicorn.err
priority=101
